# src/box/s6-rc.d/nginx/run
#!/bin/sh
# Wait for MariaDB to be ready before starting nginx
while ! mariadb-admin ping --silent 2>/dev/null; do
    sleep 1
done
exec 2>&1
exec nginx -g 'daemon off;'
