# src/box/s6-rc.d/mariadb/run
#!/usr/bin/execlineb -P
foreground { mkdir -p /run/mysqld }
foreground { chown mysql:mysql /run/mysqld }
s6-setuidgid mysql
fdmove -c 2 1
/usr/bin/mariadbd-safe
